# deploy/systemd/meadow-telegram.service
[Unit]
Description=Meadow Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@REPO_ROOT@
ExecStart=@REPO_ROOT@/.venv/bin/python -m meadow --config config/config.toml telegram
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
